# .docker/supervisor.conf
[supervisord]
nodaemon=true

[program:app]
priority=10
directory=/opt/sickbeard
command=/usr/bin/python /opt/sickbeard/SickBeard.py
user=root
autostart=true
